# scripts/nexus-sync.service
[Unit]
Description=Nexus Gmail Sync Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -s -X GET http://localhost:3000/api/cron/gmail-sync
# Add authentication header if you implemented it
# ExecStart=/usr/bin/curl -s -H "Authorization: Bearer YOUR_SECRET" -X GET http://localhost:3000/api/cron/gmail-sync

[Install]
WantedBy=multi-user.target
